# files/etc/init.d/auto_format_overlay_samba
#!/bin/sh /etc/rc.common
START=8
STOP=99

start() {
    sleep 2

    # 标记文件
    OVERLAY_MARKER="/overlay/.format_done"

    # ==============================================
    # overlay：首次/不保留配置才格式化
    # ==============================================
    if [ ! -f "$OVERLAY_MARKER" ]; then
        echo "【overlay】首次初始化 → 格式化"
        umount /dev/mmcblk0p17 2>/dev/null
        mkfs.ext4 -F -L overlay /dev/mmcblk0p17
        mkdir -p /overlay
        mount /dev/mmcblk0p17 /overlay
        touch "$OVERLAY_MARKER"
    else
        echo "【overlay】正常挂载"
        mkdir -p /overlay
        mount /dev/mmcblk0p17 /overlay 2>/dev/null
    fi

    # ==============================================
    # samba：只挂载，绝不格式化
    # ==============================================
    if [ -b "/dev/mmcblk0p18" ]; then
        mkdir -p /mnt/samba
        mount /dev/mmcblk0p18 /mnt/samba 2>/dev/null
        echo "【samba】已挂载，不格式化"
    fi
}

stop() {
    swapoff /dev/mmcblk0p19 2>/dev/null
    umount /mnt/samba 2>/dev/null
    umount /overlay 2>/dev/null
    return 0
}
