# files/etc/init.d/auto_format_overlay_samba
#!/bin/sh /etc/rc.common
START=8
STOP=99

start() {
    # 等待存储设备就绪
    sleep 3

    # 格式化标记（只格式化一次）
    OVERLAY_MARKER="/overlay/.format_done"
    SWAP_MARKER="/tmp/swap_format_done"

    # ==============================================
    # 首次启动 / sysupgrade -n 时才格式化 overlay
    # ==============================================
    if [ ! -f "$OVERLAY_MARKER" ]; then
        echo "【overlay】首次初始化，格式化 /dev/mmcblk0p17"
        
        umount /dev/mmcblk0p17 2>/dev/null
        mkfs.ext4 -F -L overlay /dev/mmcblk0p17
        mkdir -p /overlay
        mount /dev/mmcblk0p17 /overlay
        touch "$OVERLAY_MARKER"
    else
        echo "【overlay】正常挂载"
        mkdir -p /overlay
        mount /dev/mmcblk0p17 /overlay 2>/dev/null
    fi

    # ==============================================
    # 首次启动才格式化 swap（以后只启用）
    # ==============================================
    if [ ! -f "$SWAP_MARKER" ]; then
        echo "【swap】首次初始化，格式化 /dev/mmcblk0p19"
        
        swapoff /dev/mmcblk0p19 2>/dev/null
        mkswap -L swap /dev/mmcblk0p19
        swapon /dev/mmcblk0p19
        touch "$SWAP_MARKER"
    fi

    # 交给系统 fstab 统一挂载所有分区（overlay/samba/swap）
    /etc/init.d/fstab restart
}

stop() {
    swapoff /dev/mmcblk0p19 2>/dev/null
    umount /mnt/samba 2>/dev/null
    umount /overlay 2>/dev/null
    return 0
}
